# templates/default/ipfs.systemd.service.erb
[Unit]
Description=Start ipfs

[Service]
ExecStart=/usr/local/bin/ipfs daemon --migrate --routing=dhtclient
User=ipfs
Group=ipfs
Restart=always
# Raise the ulimit (max number of open files)
LimitNOFILE=<%= node['ipfs']['ulimit'] %>
# Limit the memory usage
MemoryMax=<%= node['ipfs']['memory_max'] %>

[Install]
WantedBy=multi-user.target
