# templates/default/ipfs-cluster.systemd.service.erb
[Unit]
Description=Start ipfs-cluster

[Service]
ExecStart=/usr/local/bin/ipfs-cluster-service daemon --upgrade
User=ipfs
Group=ipfs
Restart=always
# Raise the ulimit (max number of open files)
LimitNOFILE=64000

[Install]
WantedBy=multi-user.target
